# services/aiwm-worker/aiwm-worker.service
[Unit]
Description=AIWM Worker Node Client
Documentation=https://github.com/your-org/hydra-services
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User={{SERVICE_USER}}
Group={{SERVICE_USER}}
WorkingDirectory={{DEPLOY_DIR}}
EnvironmentFile={{DEPLOY_DIR}}/.env

# Start command
ExecStart=/usr/bin/node {{DEPLOY_DIR}}/main.js

# Restart policy
Restart=always
RestartSec=10
StartLimitInterval=200
StartLimitBurst=5

# Graceful shutdown
KillMode=mixed
KillSignal=SIGTERM
TimeoutStopSec=30

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=aiwm-worker

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths={{DEPLOY_DIR}} /var/run/docker.sock
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true

# Resource limits
LimitNOFILE=65536
TasksMax=4096

[Install]
WantedBy=multi-user.target
